# .root/etc/profile.d/use-xinput2.sh
#!/bin/sh

# Smooth scrolling + MacOS/Windows like touchpad pinch-to-zoom in Firefox
# https://www.reddit.com/r/linux/comments/72mfv8/psa_for_firefox_users_set_moz_use_xinput21_to/
# https://www.reddit.com/r/linux/comments/rmuh0o/finally_macoswindows_like_touchpad_zoom_gesture/
export MOZ_USE_XINPUT2=1
